# etc/worker.service
[Unit]
Description=Worker Service - Job Isolation Platform
After=network.target

[Service]
# Use single binary that auto-detects execution mode
ExecStart=/opt/worker/worker
Restart=always
RestartSec=10s

# Run as root to enable full namespace isolation capabilities
User=root
Group=root

# CRITICAL: Allow new privileges for namespace operations
NoNewPrivileges=no

# Security hardening while maintaining isolation capabilities
PrivateTmp=yes
ProtectHome=yes
ReadWritePaths=/opt/worker /sys/fs/cgroup /proc /tmp

# CRITICAL: Disable protections that block namespace operations
ProtectSystem=no
PrivateDevices=no
ProtectKernelTunables=no
ProtectControlGroups=no
RestrictRealtime=no
RestrictSUIDSGID=no
MemoryDenyWriteExecute=no

# Cgroup delegation for job resource management
Delegate=yes
DelegateControllers=cpu memory io pids
CPUAccounting=yes
MemoryAccounting=yes
IOAccounting=yes
TasksAccounting=yes
Slice=worker.slice

# Runtime settings
WorkingDirectory=/opt/worker
LimitNOFILE=4096
StandardOutput=journal
StandardError=journal

# Environment variables for the main service
Environment="WORKER_MODE=server"
Environment="LOG_LEVEL=INFO"
Environment="WORKER_CONFIG_PATH=/opt/worker/config.yml"

# Cleanup job cgroups on service stop
ExecStopPost=/bin/bash -c 'find /sys/fs/cgroup/worker.slice/worker.service -name "job-*" -type d -exec rmdir {} \; 2>/dev/null || true'

[Install]
WantedBy=multi-user.target