# debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop the service
        if systemctl is-active --quiet worker; then
            systemctl stop worker
        fi

        # Disable the service
        if systemctl is-enabled --quiet worker; then
            systemctl disable worker
        fi
        ;;
esac

exit 0